# centrallix/etc/centrallix.init.in
#!/bin/sh
#
# Startup script for the Centrallix Application Server
#
# chkconfig: 345 85 15
# processname: centrallix
# description: Centrallix Application Server
# pidfile: ##LOCALSTATEDIR##/run/centrallix.pid
# config: ##SYSCONFDIR##/centrallix.conf

# Source function library.
. /etc/rc.d/init.d/functions

centrallix=##SBINDIR##/centrallix
[ -x ${centrallix} ] || exit 0

ARGS="-q -d"
export LD_LIBRARY_PATH="$LD_LIBRARY_PATH:##LIBDIR##:##CXLIBDIR##"

[ -d ##LOCALSTATEDIR##/lock/subsys ] || mkdir -p ##LOCALSTATEDIR##/lock/subsys

start() {
    local RETVAL=0

    # Start daemon.
    printf "Starting centrallix: "
    daemon ${centrallix} $ARGS
    RETVAL=$?
    echo

    [ $RETVAL -eq 0 ] && touch ##LOCALSTATEDIR##/lock/subsys/centrallix
    return $RETVAL
}

stop() {
    local RETVAL=0

    # Stop daemons.
    printf "Stopping centrallix: "
    killproc ${centrallix}
    RETVAL=$?
    echo

    [ $RETVAL -eq 0 ] && rm -f ##LOCALSTATEDIR##/lock/subsys/centrallix
    return $RETVAL
}

# See how we were called.
case "$1" in
    start)
	start
	;;
    stop)
	stop
	;;
    status)
	status ${centrallix}
	;;
    restart)
	stop
	start
	;;
    condrestart)
	if [ -f ##LOCALSTATEDIR##/lock/subsys/centrallix ] ; then
	    stop
	    start
	fi
	;;
    *)
	printf "Usage: %s {start|stop|restart|status|condrestart}\n" "$0"
	RETVAL=1
esac

exit $RETVAL

